# MOD1-Sesion4/crontab10.crontab

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * /home/marcos/Escritorio/4.10.sh
